# config/etc/profile.d/configureXServer.sh
#!/bin/bash

export DISPLAY=:0
export LIBGL_ALWAYS_INDIRECT=0

#required for Google Chrome
export NO_AT_BRIDGE=1
