# config/etc/profile.d/configureDocker.sh
#!/usr/bin/env bash

# connect to Windows Docker
export DOCKER_HOST=tcp://localhost:2375
